# mateng/lib/udev/rules.d/60-vlan-network-interface.rules
ACTION=="add", SUBSYSTEM=="net", RUN+="vlan-network-interface UDEV"
