# mateng/var/lib/dpkg/info/lxd.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/lxd
/etc/sysctl.d
/etc/sysctl.d/10-lxd-inotify.conf
/etc/apparmor.d
/etc/apparmor.d/usr.lib.lxd.lxd-bridge-proxy
/etc/init.d
/etc/init.d/lxd
/etc/init
/etc/init/lxd.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lxd.socket
/lib/systemd/system/lxd.service
/lib/systemd/system/lxd-containers.service
/lib/systemd/system/lxd-bridge.service
/var
/var/lib
/var/lib/lxd
/var/log
/var/log/lxd
/usr
/usr/lib
/usr/lib/lxd
/usr/lib/lxd/lxd-bridge.start
/usr/lib/lxd/lxd-bridge-proxy
/usr/lib/lxd/profile-config
/usr/lib/lxd/shutdown
/usr/lib/lxd/lxd-bridge
/usr/bin
/usr/bin/lxd
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/lxd
/usr/share/doc
/usr/share/doc/lxd
/usr/share/doc/lxd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/lxd.1.gz
/usr/share/doc/lxd/changelog.Debian.gz
