# mateng/var/lib/dpkg/info/open-iscsi.conffiles
/etc/default/open-iscsi
/etc/init.d/iscsid
/etc/init.d/open-iscsi
/etc/iscsi/iscsid.conf
