# mateng/var/lib/dpkg/info/open-iscsi.list
/.
/etc
/etc/iscsi
/etc/iscsi/iscsid.conf
/etc/init.d
/etc/init.d/iscsid
/etc/init.d/open-iscsi
/etc/network
/etc/network/if-up.d
/etc/default
/etc/default/open-iscsi
/sbin
/sbin/iscsid
/sbin/iscsi-iname
/sbin/iscsiadm
/sbin/iscsistart
/sbin/iscsi_discovery
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/open-iscsi.service
/lib/systemd/system/iscsid.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/70-iscsi-network-interface.rules
/lib/open-iscsi
/lib/open-iscsi/startup-checks.sh
/lib/open-iscsi/net-interface-handler
/lib/open-iscsi/activate-storage.sh
/lib/open-iscsi/logout-all.sh
/lib/open-iscsi/umountiscsi.sh
/lib/modules-load.d
/lib/modules-load.d/open-iscsi.conf
/usr
/usr/bin
/usr/share
/usr/share/doc
/usr/share/doc/open-iscsi
/usr/share/doc/open-iscsi/NEWS.Debian.gz
/usr/share/doc/open-iscsi/README.Debian.gz
/usr/share/doc/open-iscsi/changelog.Debian.gz
/usr/share/doc/open-iscsi/README.gz
/usr/share/doc/open-iscsi/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/iscsid.8.gz
/usr/share/man/man8/iscsi-iname.8.gz
/usr/share/man/man8/iscsistart.8.gz
/usr/share/man/man8/iscsiadm.8.gz
/usr/share/man/man8/iscsi_discovery.8.gz
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/iscsi
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/local-top
/usr/share/initramfs-tools/scripts/local-top/iscsi
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/open-iscsi
/usr/bin/iscsiadm
