# mateng/var/lib/dpkg/info/popularity-contest.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/popularity-contest
/etc/cron.d
/usr
/usr/sbin
/usr/sbin/popularity-contest
/usr/sbin/popcon-largest-unused
/usr/share
/usr/share/popularity-contest
/usr/share/popularity-contest/debian-popcon.gpg
/usr/share/popularity-contest/default.conf
/usr/share/popularity-contest/popcon-upload
/usr/share/doc
/usr/share/doc/popularity-contest
/usr/share/doc/popularity-contest/examples
/usr/share/doc/popularity-contest/examples/bin
/usr/share/doc/popularity-contest/examples/bin/popcon-process.sh
/usr/share/doc/popularity-contest/examples/bin/popanal.py
/usr/share/doc/popularity-contest/examples/bin/README.examples
/usr/share/doc/popularity-contest/examples/bin/popcon.pl
/usr/share/doc/popularity-contest/examples/bin/prepop.pl
/usr/share/doc/popularity-contest/examples/bin/gensections.pl
/usr/share/doc/popularity-contest/examples/bin/popcon-stat.pl
/usr/share/doc/popularity-contest/examples/cgi-bin
/usr/share/doc/popularity-contest/examples/cgi-bin/popcon.cgi
/usr/share/doc/popularity-contest/changelog.gz
/usr/share/doc/popularity-contest/FAQ.gz
/usr/share/doc/popularity-contest/copyright
/usr/share/doc/popularity-contest/README.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/popularity-contest.8.gz
/usr/share/man/man8/popcon-largest-unused.8.gz
