# mateng/var/lib/dpkg/info/preload.list
/.
/var
/var/log
/var/lib
/var/lib/preload
/usr
/usr/sbin
/usr/sbin/preload
/usr/share
/usr/share/doc
/usr/share/doc/preload
/usr/share/doc/preload/TODO
/usr/share/doc/preload/THANKS
/usr/share/doc/preload/changelog.Debian.gz
/usr/share/doc/preload/proposal.txt.gz
/usr/share/doc/preload/AUTHORS
/usr/share/doc/preload/copyright
/usr/share/doc/preload/README.gz
/usr/share/doc/preload/NEWS.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/preload.8.gz
/etc
/etc/init.d
/etc/init.d/preload
/etc/preload.conf
/etc/logrotate.d
/etc/logrotate.d/preload
/etc/default
/etc/default/preload
