# mateng/var/lib/dpkg/info/preload.conffiles
/etc/init.d/preload
/etc/preload.conf
/etc/logrotate.d/preload
/etc/default/preload
