# mateng/var/lib/dpkg/info/screen.list
/.
/etc
/etc/tmpfiles.d
/etc/screenrc
/etc/init.d
/etc/init.d/screen-cleanup
/usr
/usr/bin
/usr/bin/screen
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/screen-cleanup.conf
/usr/share
/usr/share/info
/usr/share/info/screen.info-5.gz
/usr/share/info/screen.info-3.gz
/usr/share/info/screen.info-2.gz
/usr/share/info/screen.info-6.gz
/usr/share/info/screen.info-1.gz
/usr/share/info/screen.info-4.gz
/usr/share/info/screen.info.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/screen
/usr/share/screen
/usr/share/screen/utf8encodings
/usr/share/screen/utf8encodings/bf
/usr/share/screen/utf8encodings/01
/usr/share/screen/utf8encodings/03
/usr/share/screen/utf8encodings/18
/usr/share/screen/utf8encodings/c2
/usr/share/screen/utf8encodings/cc
/usr/share/screen/utf8encodings/c3
/usr/share/screen/utf8encodings/04
/usr/share/screen/utf8encodings/02
/usr/share/screen/utf8encodings/cd
/usr/share/screen/utf8encodings/c6
/usr/share/screen/utf8encodings/19
/usr/share/screen/utf8encodings/c8
/usr/share/screen/utf8encodings/c4
/usr/share/screen/utf8encodings/d6
/usr/share/screen/utf8encodings/a1
/usr/share/screen/utf8encodings/c7
/usr/share/doc
/usr/share/doc/screen
/usr/share/doc/screen/NEWS.3.7.gz
/usr/share/doc/screen/NEWS.gz
/usr/share/doc/screen/examples
/usr/share/doc/screen/examples/screenrc
/usr/share/doc/screen/examples/etcscreenrc
/usr/share/doc/screen/NEWS.3.5.gz
/usr/share/doc/screen/FAQ.gz
/usr/share/doc/screen/copyright
/usr/share/doc/screen/changelog.Debian.gz
/usr/share/doc/screen/TODO
/usr/share/doc/screen/terminfo
/usr/share/doc/screen/terminfo/screeninfo.src
/usr/share/doc/screen/terminfo/README.terminfo
/usr/share/doc/screen/terminfo/checktc.c.gz
/usr/share/doc/screen/terminfo/8bits
/usr/share/doc/screen/terminfo/tetris.c
/usr/share/doc/screen/terminfo/test.txt.gz
/usr/share/doc/screen/terminfo/README
/usr/share/doc/screen/terminfo/screencap
/usr/share/doc/screen/fdpat.ps.gz
/usr/share/doc/screen/window_to_display.ps.gz
/usr/share/doc/screen/README
/usr/share/doc/screen/NEWS.Debian.gz
/usr/share/doc/screen/README.Debian.gz
/usr/share/doc/screen/NEWS.3.6.gz
/usr/share/doc/screen/NEWS.3.9.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/screen.1.gz
