# mateng/var/lib/dpkg/info/screen.conffiles
/etc/init.d/screen-cleanup
/etc/screenrc
