# mateng/var/lib/dpkg/info/snapd.list
/.
/usr
/usr/lib
/usr/lib/snapd
/usr/lib/snapd/snap-gdb-shim
/usr/lib/snapd/snap-mgmt
/usr/lib/snapd/snap-exec
/usr/lib/snapd/snap-device-helper
/usr/lib/snapd/snap-confine
/usr/lib/snapd/complete.sh
/usr/lib/snapd/snapd
/usr/lib/snapd/snap-repair
/usr/lib/snapd/system-shutdown
/usr/lib/snapd/info
/usr/lib/snapd/snapd.core-fixup.sh
/usr/lib/snapd/snap-seccomp
/usr/lib/snapd/snap-discard-ns
/usr/lib/snapd/etelpmoc.sh
/usr/lib/snapd/snap-update-ns
/usr/share
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/io.snapcraft.snapd.policy
/usr/share/doc
/usr/share/doc/snapd
/usr/share/doc/snapd/changelog.gz
/usr/share/doc/snapd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/snap-confine.1.gz
/usr/share/man/man1/snap.1.gz
/usr/share/man/man5
/usr/share/man/man5/snap-discard-ns.5.gz
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/snap
/usr/share/dbus-1
/usr/share/dbus-1/services
/usr/share/dbus-1/services/io.snapcraft.Launcher.service
/usr/share/dbus-1/services/io.snapcraft.Settings.service
/usr/bin
/usr/bin/snapctl
/usr/bin/snap
/usr/bin/snapfuse
/etc
/etc/apparmor.d
/etc/apparmor.d/usr.lib.snapd.snap-confine.real
/etc/profile.d
/etc/profile.d/apps-bin-path.sh
/lib
/lib/systemd
/lib/systemd/system-generators
/lib/systemd/system-generators/snapd-generator
/lib/systemd/system
/lib/systemd/system/snapd.service
/lib/systemd/system/snapd.snap-repair.timer
/lib/systemd/system/snapd.socket
/lib/systemd/system/snapd.core-fixup.service
/lib/systemd/system/snapd.snap-repair.service
/lib/systemd/system/snapd.autoimport.service
/lib/systemd/system/snapd.system-shutdown.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/66-snapd-autoimport.rules
/snap
/var
/var/cache
/var/cache/snapd
/var/lib
/var/lib/snapd
/var/lib/snapd/auto-import
/var/lib/snapd/environment
/var/lib/snapd/lib
/var/lib/snapd/lib/gl32
/var/lib/snapd/lib/gl
/var/lib/snapd/lib/vulkan
/var/lib/snapd/desktop
/var/lib/snapd/snaps
/var/lib/snapd/snaps/partial
/var/lib/snapd/void
/var/lib/snapd/apparmor
/var/lib/snapd/apparmor/snap-confine
/var/lib/snapd/firstboot
/var/snap
/usr/bin/ubuntu-core-launcher
/lib/udev/snappy-app-dev
/etc/X11/Xsession.d/65snappy
