# mateng/etc/init/zram-config.conf
description	"Initializes zram swaping"
author		"Adam Conrad <adconrad@canonical.com>"

start on runlevel [2345]

pre-start exec /usr/bin/init-zram-swapping

pre-stop exec /usr/bin/end-zram-swapping

